# docker/php-fpm/supervisord.d/php-fpm.conf
[program:php-fpm]
command=php-fpm
user=root
